# test/cmocka/supervisord.conf
[supervisord]
nodaemon=true

[program:ubusd]
command=/bin/bash -c "/usr/sbin/ubusd"

[program:rpcd]
command=/bin/bash -c "/usr/sbin/rpcd"

[program:sshd]
command=/usr/sbin/sshd -D
